# conf/systemd.timer.j2
[Unit]
Description=Run backup {{ app }} regularly

[Timer]
OnCalendar={{ on_calendar }}

[Install]
WantedBy=timers.target
