# conf/systemd-cleanup-logs.timer
[Unit]
Description=Run __APP__ log cleanup regularly

[Timer]
OnCalendar=Weekly

[Install]
WantedBy=timers.target
